# cloudflare-speedtest.service
# ~/.config/systemd/user/cloudflare-speedtest.service
[Unit]
Description=Cloudflare Speed Test (background)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/andrei/.local/bin/cloudflare-speed-cli-bg
RemainAfterExit=yes
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/home/andrei/.local/bin

[Install]
WantedBy=default.target
